# aware-updater.service
[Unit]
Description=AWAre Updater Service

[Service]
ExecStart=/home/iata/aware-updater/updater
Restart=always

[Install]
WantedBy=default.target
